# templates/reference-projects/rpm-service-app/tizen/rpm/sources/hello-service.service
[Unit]
Description=hello-service reference daemon
After=network.target

[Service]
Type=simple
EnvironmentFile=-/etc/hello-service.env
ExecStart=/usr/bin/hello-service
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
